# ubuntu20mariadb/src/supervisord.conf
[supervisord]
nodaemon=true

[program:mariaserver]
command=/init.sh