# ubuntu20apache2php74/src/supervisord.conf
[supervisord]
nodaemon=true

[program:php74fpm]
command=service php7.4-fpm start

[program:apache2]
command=apache2ctl -D FOREGROUND
