# ubuntu20apache2php81/src/supervisord.conf
[supervisord]
nodaemon=true

[program:php81fpm]
command=service php8.1-fpm start

[program:apache2]
command=apache2ctl -D FOREGROUND
